# crontab
*/1 * * * * root /bin/bash /tmp/transport.sh >> /var/log/cron 2>&1

